# deploy/discord-mcp.service
# systemd unit — Discord MCP (service HTTP 24/7).
# ⚠️ Installer dans /etc/systemd/system/discord-mcp.service sur le VPS Allemagne (100.64.0.1).
# ⚠️ EnvironmentFile = perms 600, propriété root. Contient DISCORD_MCP_HTTP_TOKEN + DISCORD_SECRETS_PATH.
[Unit]
Description=Discord MCP multi-bot — transport HTTP
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/discord-mcp
EnvironmentFile=/etc/discord-mcp/env
ExecStart=/usr/bin/node http.js
Restart=always
RestartSec=3

# Durcissement (le service ne lit que son dossier + son secret).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/discord-mcp
# bind 127.0.0.1 only (cf http.js) → pas d'écoute publique. Accès via Tailscale.

[Install]
WantedBy=multi-user.target
